# etc/rc.local
# Print the IP address
_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

/home/pi/bin/keepalive.sh
/home/pi/bin/run.sh &

exit 0
