# etc/rc.local
/home/pi/bin/keepalive.sh &
/home/pi/bin/run.sh &

exit 0
